# deploy/micilla-eduweb.service
# systemd unit for MiCilla EduWeb on a VPS.
#
# Install:
#   sudo cp deploy/micilla-eduweb.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now micilla-eduweb
#
# Adjust User, WorkingDirectory, and EnvironmentFile to match your server.

[Unit]
Description=MiCilla EduWeb Next.js application
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/micilla-eduweb
Environment=NODE_ENV=production
EnvironmentFile=/var/www/micilla-eduweb/.env
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
